# other/fido.service
[Unit]
Description=overmind for fido_jr

[Service]
ExecStart=/home/root/service/start.sh

[Install]
WantedBy=multi-user.target
